# systemd/job-automation-orchestrator.service
[Unit]
Description=Job Automation Orchestrator Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/automation
ExecStart=/usr/bin/node /home/ubuntu/automation/apps/orchestrator/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/automation/.env
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=job-orchestrator

[Install]
WantedBy=multi-user.target
