# scripts/candidate-backup.timer
[Unit]
Description=Candidate Knowledge Automated Daily Backup Schedule (3:00 AM UTC / 8:30 AM IST)

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
